# oya_package/files/etc/uci-defaults/oya_01_init
#!bin/sh

#Give permission for init.d scrips for Oya
chmod +x /etc/init.d/oya

exit 0
