# oya_package/files/etc/uci-defaults/oya_inits
#!bin/sh

chmod +x /etc/init.d/oya
exit 0
