# oya_package/files/etc/uci-defaults/z_oya_01_init
#!bin/sh

#Give permission for init.d scrips for Oya
chmod +x /etc/init.d/oya
/etc/init.d/oya start
chmod +x /usr/bin/device_id
exit 0
