# oya_package/files/etc/init.d/oya
#!/bin/sh /etc/rc.common

START=99

start()
{
	sh /Oya/main.sh
}
stop() {
	#kill the oya processes
}
restart(){
	#restart oya processes
}
boot() {
	echo "Boot function runs" > /root/start.txt
}
