# oya_package/files/etc/init.d/oya
#!/bin/sh /etc/rc.common

START=10

start()
{
	sh /Oya/main.sh
}
boot() {
  sh /Oya/main.sh
}
