# oya_package/files/etc/init.d/oya
#!/bin/sh /etc/rc.common

START=99
STOP=99

start() {
	sh /Oya/main.sh
}
boot() {
	echo "Boot function runs" > /root/start.txt
}
