# oya_package/files/etc/init.d/recording
#!/bin/sh /etc/rc.common

START=99
STOP=99

start() {
  if [ ! -f /Oya/pid/recording ]; then
    sh /Oya/recording.sh
  fi
}

stop() {
  recPid=`recording pid get`
  arecPid=`recording apid get`
  kill ${recPid}
  kill ${arecPid}
  recording pid remove
  recording apid remove
}

boot() {
	/etc/init.d/recording enable
}
